# deploy/mydoorcamera-capture.service
[Unit]
Description=mydoorcamera capture daemon (Rust: activity-gated hardware H.264 recorder)
After=multi-user.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/mydoorcamera/recorder
ExecStart=/root/mydoorcamera/recorder/target/release/recorder
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
